# llupdatetool/Library/LaunchAgents/edu.mit.BitBarLaunch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>edu.mit.BitBarLaunch</string>
	<key>ProgramArguments</key>
	<array>
		<string>sh</string>
		<string>/usr/local/updateTool/openBitBar.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>399900</integer>
</dict>
</plist>
